# php/crontab
* * * * * echo 'crontab test' >> /usr/cron/cron.log 
