# php/supervisord.conf
[supervisord]
nodaemon=true

[program:php]
command=php-fpm

[program:cron]
command=service cron start
